# memguard.service
[Unit]
Description=memguard desktop memory pressure daemon
After=systemd-logind.service dbus.service
Wants=dbus.service

[Service]
Type=simple
ExecStart=/usr/bin/memguard
Restart=on-failure
RestartSec=5
WatchdogSec=30

[Install]
WantedBy=multi-user.target
