# system-tune/memguard-system-tune.service
[Unit]
Description=One-time memguard system tuning
[Service]
Type=oneshot
ExecStart=/usr/bin/memguard-system-tune
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
